# integration/packer/alluxio-presto-sandbox-ami/supervisor/supervisord.conf
[supervisord]
childlogdir=/var/log/supervisor/
logfile=/var/log/supervisor/supervisord.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]

[inet_http_server]
port=9001

[program:mariadb]
command=/usr/libexec/mysqld --basedir=/usr --datadir=/var/lib/mysql --plugin-dir=/usr/lib64/mysql/plugin --user=mysql --log-error=/var/log/mariadb/mariadb.log --pid-file=/var/run/mariadb/mariadb.pid --socket=/var/lib/mysql/mysql.sock
priority=1
startsecs=2

[program:hadoop_namenode]
command=bash -c '. ~/.bashrc; exec ./bin/hdfs namenode'
directory=/opt/hadoop
autostart=true
environment=HOME="/home/ec2-user",USER="ec2-user"
priority=2
startsecs=5
user=ec2-user

[program:hadoop_snamenode]
command=bash -c '. ~/.bashrc; exec ./bin/hdfs secondarynamenode'
directory=/opt/hadoop
autostart=true
environment=HOME="/home/ec2-user",USER="ec2-user"
priority=2
startsecs=5
user=ec2-user

[program:hadoop_datanode]
command=bash -c '. ~/.bashrc; exec ./bin/hdfs datanode'
directory=/opt/hadoop
autostart=true
environment=HOME="/home/ec2-user",USER="ec2-user"
priority=2
startsecs=5
user=ec2-user

[program:hive_metastore]
command=bash -c '. ~/.bashrc; exec ./bin/hive --service metastore'
directory=/opt/hive
autostart=true
environment=HOME="/home/ec2-user",USER="ec2-user"
priority=3
startsecs=5
user=ec2-user

[program:alluxio_master]
command=bash -c '. ~/.bashrc; ./bin/launch-process master -c'
directory=/opt/alluxio
autostart=true
environment=HOME="/home/ec2-user",USER="ec2-user"
priority=4
startsecs=3
user=ec2-user

[program:alluxio_job_master]
command=bash -c '. ~/.bashrc; ./bin/launch-process job_master -c'
directory=/opt/alluxio
autostart=true
environment=HOME="/home/ec2-user",USER="ec2-user"
priority=4
startsecs=3
user=ec2-user

[program:alluxio_job_worker]
command=bash -c '. ~/.bashrc; ./bin/launch-process job_worker -c'
directory=/opt/alluxio
autostart=true
environment=HOME="/home/ec2-user",USER="ec2-user"
priority=6
startsecs=3
user=ec2-user

[program:alluxio_proxy]
command=bash -c '. ~/.bashrc; ./bin/launch-process proxy -c'
directory=/opt/alluxio
autostart=true
environment=HOME="/home/ec2-user",USER="ec2-user"
priority=6
startsecs=3
user=ec2-user

[program:alluxio_worker]
command=bash -c '. ~/.bashrc; ./bin/alluxio-mount.sh SudoMount && ./bin/launch-process worker -c'
directory=/opt/alluxio
autostart=true
environment=HOME="/home/ec2-user",USER="ec2-user"
priority=6
startsecs=3
user=ec2-user

[program:presto]
command=bash -c '. ~/.bashrc; /opt/presto/bin/launcher run'
autostart=true
environment=HOME="/home/ec2-user",USER="ec2-user"
priority=5
startsecs=5
user=ec2-user
